# systemd/prestamos-backup.timer
[Unit]
Description=Run Prestamos SQLite backup daily

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
Unit=prestamos-backup.service

[Install]
WantedBy=timers.target
